# scdp.cron
#
# cron job for the scdp package
#
* * * * *	root	/usr/bin/scdp.cron
